# backend/lukete/lukete.service
[Unit]
Description=Lukete Refund Dashboard (Streamlit)
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/fde1/lukete
Environment="PATH=/home/ec2-user/fde1/lukete/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=-/home/ec2-user/fde1/lukete/.env
ExecStart=/home/ec2-user/fde1/lukete/venv/bin/streamlit run dashboard.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
